# src/containers/host/supervisord.conf
[supervisord]
nodaemon=true
user=klangk
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# klangkd binds the UDS and owns the proxy (nginx) as its child (#1396). The old
# two-program layout (separate [program:nginx]) is collapsed into this single
# entry; proxy crash-restart now happens inside klangkd's lifespan watchdog.
# Supervisord remains the container's PID-1 (restarts klangkd on crash, reaps
# zombies); it no longer manages the proxy directly.
[program:uvicorn]
command=python3 -m klangk.launcher --config=none
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
stopasgroup=true
killasgroup=true
